# managed/systemd/openclaw-gateway.service.d/acp-harness.conf
[Service]
EnvironmentFile=%h/.openclaw/acp-harness.env
